# tlsd.service
[Service]
ExecStart=/usr/bin/tlsd
User=daemon

[Install]
WantedBy=default.target
