# doc/tlsd.service
[Service]
ExecStart=/usr/bin/tlsd

[Install]
WantedBy=default.target
