# doc/tlsd.socket
[Unit]
Description=TLSD Server Activation Socket

[Socket]
ListenStream=/run/tlsd.socket
RemoveOnStop=true

[Install]
WantedBy=sockets.target
